# scipts/rsyslog/crontab
* * * * * /usr/sbin/logrotate -f /etc/logrotate.d/bash-history
